# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/etc/init.d/httpd start
autostart=true
autorestart=true